# nodepool/scripts/perf-iostat.conf
# gather perf top log
#

description     "gather perf iostat"

start on start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

pre-start script
        test -d /var/log/perf || mkdir -p -m0755 /var/log/perf/
        test -e /usr/bin/iostat || apt-get -y --force-yes install sysstat
end script

exec iostat -t -x 5 >>/home/jenkins/workspace/dsvm-tempest-nova-network/logs/perf-iostat.txt

